# distribution/deb/etc/service/orkaudio/run
#!/bin/bash
set -x
exec 2>&1

# exec chpst -u USER COMMAND
exec /usr/sbin/orkaudio
