# distribution/deb/orkaudio/etc/service/orkaudio/run
#!/bin/bash
exec 2>&1

exit_script() {
    echo "Shutdown!"
    CHILD_PID=$(cat /var/log/orkaudio/orkaudio.lock | head -n1)

    echo "Sending Kill to $CHILD_PID"
    trap - SIGINT SIGTERM # clear the trap
    kill -- $CHILD_PID # Sends SIGTERM to child/sub processes
    rm /var/log/orkaudio/orkaudio.lock
}

anywait(){
    for pid in "$@"; do
        while kill -0 "$pid"; do
            sleep 2
        done
    done
}

ulimit -n 65535
ulimit -c unlimited

trap exit_script SIGINT SIGTERM

# exec chpst -u USER COMMAND
export LIBC_FATAL_STDERR_=1
epoch=$(date +%s)
SEGFAULT_SIGNALS=all  SEGFAULT_OUTPUT_NAME="/var/log/orkaudio/orkaudio.$epoch.trace" LD_PRELOAD=/lib/x86_64-linux-gnu/libSegFault.so /usr/sbin/orkaudio

#wait 1 sec for child lock to be created.
sleep 1

CHILD_PID=$(cat /var/log/orkaudio/orkaudio.lock | head -n1 )

#watch on pid
echo "Starting Watch on $CHILD_PID"
while kill -0 $CHILD_PID 2> /dev/null; do sleep 1; done;
